# supervisor/mariadb_backup.ini
[program:mariadb_backup_service]
command = /usr/local/lib/maridb-backup-service/.venv/bin/python /usr/local/lib/maridb-backup-service/run.py
redirect_stderr = True
stdout_logfile = var/log/mariadb-backup-service/mariadb-service.log
autorestart = True
autostart = True

